# src/tally.service.in
[Unit]
Description=Voctomix Tally Light
Documentation=https://c3voc.de/wiki/projects:tallycom

[Service]
EnvironmentFile=/etc/tally.conf
ExecStart=/bin/sh -c '@prefix@/bin/tallypi-tally ${host} ${source}'
Restart=always
RestartSec=5
ProtectHome=yes
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
